# s6-overlay/s6-rc.d/bot/run
#!/command/with-contenv sh
# Load environment and run the relay bot
cd /relay
exec python3 bot.py
