# s6-overlay-session/s6-rc.d/claude-update/up
#!/command/execlineb -P
foreground { /bin/sh -lc "if [ \"\${CLAUDE_AUTO_UPDATE:-0}\" = \"1\" ]; then claude update </dev/null >/dev/null 2>&1 || true; fi" }
ifelse { true } { }
