# s6-overlay/s6-rc.d/queue-backup/run
#!/command/execlineb -P
/bin/bash /relay/scripts/queue-backup-loop.sh
