# s6-overlay-session/s6-rc.d/codex-bot/run
#!/command/with-contenv sh
# Standalone Codex Telegram bot — only runs in codex session containers
# Requires CODEX_BOT_TOKEN to be set

# Skip if not a codex session or no bot token
if [ "${SESSION_TYPE}" != "codex" ] || [ -z "${CODEX_BOT_TOKEN}" ]; then
  echo "[codex-bot] Skipping — SESSION_TYPE=${SESSION_TYPE}, CODEX_BOT_TOKEN=${CODEX_BOT_TOKEN:+set}"
  exec sleep infinity
fi

echo "[codex-bot] Starting standalone bot for session ${SESSION_NAME}"
exec node /root/relay/scripts/codex-bot.js
