# s6-overlay-session/s6-rc.d/codex-bot/run
#!/command/with-contenv sh
# Standalone Codex Telegram bot — only runs in codex session containers
# Requires CODEX_BOT_TOKEN to be set
# DISABLED: webhook mode handles all Telegram updates via relay-api now.
# codex-bot.js polling is redundant and can conflict with webhooks.

echo "[codex-bot] Disabled — webhook mode active. codex-bot.js polling no longer needed."
exec sleep infinity
