# s6-overlay/s6-rc.d/auto-restart/run
#!/command/execlineb -P
/bin/bash /relay/scripts/auto-restart-loop.sh
