# s6-overlay/s6-rc.d/claude-update/up
#!/command/execlineb -P
foreground { claude update --yes }
ifelse { true } { }
