# s6-overlay/s6-rc.d/mcp-server/run
#!/command/with-contenv sh
# MCP server service — supervised by s6
# Env vars SESSION_NAME and TELEGRAM_THREAD_ID are injected via docker-compose env_file / environment
exec /relay/scripts/mcp-server-wrapper.sh
